# serverKent/system/llamacpp-manager.service
[Unit]
Description=Llama.cpp Manager Service (sample)
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=5
User=ai-service
Group=ai-service
WorkingDirectory=/opt/llamacpp
ExecStart=/usr/bin/python3 /opt/llamacpp/llama_manager.py
Environment=LLAMA_CONFIG=/etc/llamacpp/config.json
Environment=PYTHONPATH=/opt/llamacpp

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/llamacpp

[Install]
WantedBy=multi-user.target

